# .config/systemd/user/vnc_server_x0.timer
[Unit]
Description=NC Server - X0

[Timer]
OnBootSec=30s

[Install]
WantedBy=timers.target
