# .config/autostart/imwheel.desktop
[Desktop Entry]
Exec=/usr/bin/imwheel
Icon=
Name=imwheel
Path=
Terminal=False
Type=Application
